# deploy/systemd/kytickar-frontend.service
[Unit]
Description=Kytickar Frontend (Vite Preview)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kytickar/app
ExecStart=/bin/bash -c 'npm run preview -- --host 0.0.0.0 --port 5173'
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
